# aurora-win/scripts/schedule_audit_weekly.xml
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
<Triggers>
<CalendarTrigger>
<WeeklySchedule>
<DaysOfWeek><Sunday/></DaysOfWeek>
<WeeksInterval>1</WeeksInterval>
</WeeklySchedule>
<StartBoundary>2025-11-09T03:00:00</StartBoundary>
</CalendarTrigger>
</Triggers>
<Principals><Principal id="Author" RunLevel="LeastPrivilege"/></Principals>
<Settings><MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy></Settings>
<Actions>
<!-- [FIX] 'scripts/' 경로 추가 -->
<Exec><Command>python</Command><Arguments>scripts/audit_verify.py data/audit.log</Arguments></Exec>
<!-- [FIX] 'scripts/' 경로 추가 -->
<Exec><Command>python</Command><Arguments>scripts/audit_compact_weekly.py --in data/audit.log --out data/archive</Arguments></Exec>
</Actions>
</Task>